# com.stevew.commute-tracker-fetcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.stevew.commute-tracker-fetcher</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>REPLACE_WITH_ABSOLUTE_PATH/fetcher.py</string>
  </array>

  <!-- Run every day at 22:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <!-- Also run when the Mac boots / logs in — catches up anything missed -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Working dir so the script can run git commands in the repo -->
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_ABSOLUTE_PATH</string>

  <!-- Logs -->
  <key>StandardOutPath</key>
  <string>REPLACE_WITH_ABSOLUTE_PATH/fetcher.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_ABSOLUTE_PATH/fetcher.log</string>

  <!-- Ensure PATH includes git when launched by launchd -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
